# packages/set_oxygen_theme/root/etc/systemd/system/set_oxygen_theme.service
[Unit]
Description=Service that installs the script to autostart to enable oxygen theme and then removes itsself
After=display-manager.service

[Service]
ExecStart=/usr/local/bin/set_oxygen_theme.elf
User=root
Group=root
	

[Install]
WantedBy=multi-user.target
